# examples/deployments/personal-http/com.kb-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  ~/Library/LaunchAgents/com.kb-mcp.plist

  macOS launchd LaunchAgent for the personal-http deployment. A LaunchAgent
  (vs a LaunchDaemon) runs in your user session, which matches the
  personal-machine assumption: one user, no privileged operations.

  Install:
      mkdir -p ~/Library/LaunchAgents
      cp com.kb-mcp.plist ~/Library/LaunchAgents/com.kb-mcp.plist
      # edit the ProgramArguments path and WorkingDirectory before loading
      launchctl load ~/Library/LaunchAgents/com.kb-mcp.plist

  Restart:
      launchctl unload ~/Library/LaunchAgents/com.kb-mcp.plist
      launchctl load   ~/Library/LaunchAgents/com.kb-mcp.plist

  Tail logs (the StandardOutPath/StandardErrorPath below):
      tail -f /tmp/kb-mcp.out /tmp/kb-mcp.err

  Note on session lifetime: LaunchAgents are unloaded when you log out. To
  keep kb-mcp running across logout, switch to a LaunchDaemon at
  /Library/LaunchDaemons/ — but that needs sudo and a system user, which is
  past what this recipe covers. Stick with LaunchAgent unless you have a
  reason.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kb-mcp</string>

    <!-- Replace /Users/you/.cargo/bin/kb-mcp with the absolute path to your
         kb-mcp binary. `serve` reads transport / kb_path / model from the
         kb-mcp.toml in WorkingDirectory below. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/you/.cargo/bin/kb-mcp</string>
        <string>serve</string>
    </array>

    <!-- CWD = directory holding kb-mcp.toml so Config::discover() picks
         it up (priority 2 — CWD). Replace with your actual path. -->
    <key>WorkingDirectory</key>
    <string>/Users/you/kb-mcp</string>

    <!-- Start at login + restart on crash. The combination keeps the
         daemon warm across editor restarts. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Bump to debug while troubleshooting. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>

    <!-- Logs land here. Rotate manually or wire your own logrotate / newsyslog. -->
    <key>StandardOutPath</key>
    <string>/tmp/kb-mcp.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/kb-mcp.err</string>
</dict>
</plist>
