# examples/deployments/personal-http/kb-mcp.user.service
# ~/.config/systemd/user/kb-mcp.service
#
# systemd USER unit (not system-level) for the personal-http deployment.
# A user unit lives under your account and runs without root, which matches
# the personal-machine assumption: one user, one KB, no multi-tenancy.
#
# Enable with:
#     mkdir -p ~/.config/systemd/user/
#     cp kb-mcp.user.service ~/.config/systemd/user/kb-mcp.service
#     # edit ExecStart / Environment to point at your binary + kb-mcp.toml
#     systemctl --user daemon-reload
#     systemctl --user enable --now kb-mcp.service
#
# Optional but recommended: keep the unit running after logout so the daemon
# stays warm across login sessions.
#     sudo loginctl enable-linger "$USER"

[Unit]
Description=kb-mcp loopback HTTP MCP server (personal)
After=default.target

[Service]
Type=simple

# CWD = directory holding kb-mcp.toml so Config::discover() picks it up
# (priority 2 — CWD). Replace %h/kb-mcp with wherever you keep the toml.
WorkingDirectory=%h/kb-mcp

# Replace with the absolute path to your kb-mcp binary if it is not on $PATH
# inside the systemd user environment (cargo-installed binaries usually are
# at ~/.cargo/bin which user units inherit on most distros).
ExecStart=%h/.cargo/bin/kb-mcp serve

# Keep the daemon warm — restart on transient failures, don't burn cycles.
Restart=on-failure
RestartSec=5s

# RUST_LOG controls verbosity. info gives the "loaded config source=..."
# startup line plus warnings; bump to debug when troubleshooting.
Environment=RUST_LOG=info

# stdout/stderr → user journal. Tail with:
#     journalctl --user -u kb-mcp.service -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
