# examples/deployments/personal-http/kb-mcp-task.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Windows Task Scheduler XML for the personal-http deployment.

  Imports as a single task that starts kb-mcp at logon and restarts on
  failure. The personal-machine assumption (one user, loopback only)
  means we DO NOT run with elevated privileges and DO NOT register as
  a service via the Service Control Manager — a per-user scheduled
  task is the right granularity.

  Install (PowerShell, no admin elevation needed for AT_LOGON):
      schtasks /Create /TN "kb-mcp" /XML kb-mcp-task.xml

  Verify:
      schtasks /Query /TN "kb-mcp" /V /FO LIST

  Start now (without waiting for next logon):
      schtasks /Run /TN "kb-mcp"

  Stop:
      schtasks /End /TN "kb-mcp"

  Uninstall:
      schtasks /Delete /TN "kb-mcp" /F

  Edit BEFORE importing:
  - <Command> path: replace with the absolute path to your kb-mcp.exe.
  - <WorkingDirectory>: replace with the directory holding kb-mcp.toml so
    Config::discover() picks it up via CWD (priority 2).

  An alternative for a Windows install is `nssm` (https://nssm.cc/) which
  registers kb-mcp as a true Windows service rather than a scheduled task.
  Use nssm if you need the daemon running before any user logs in. For
  single-user development workstations, AT_LOGON is simpler and avoids
  needing admin rights.
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>kb-mcp loopback HTTP MCP server (personal)</Description>
    <URI>\kb-mcp</URI>
  </RegistrationInfo>

  <Triggers>
    <LogonTrigger>
      <Enabled>true</Enabled>
    </LogonTrigger>
  </Triggers>

  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>

  <Settings>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <Hidden>false</Hidden>
    <RestartOnFailure>
      <!--
        Windows Task Scheduler requires Interval >= PT1M (1 minute).
        Earlier values like PT5S are accepted by the XML schema but rejected
        at registration time with:
          "(66,24):Interval:PT5S — value not allowed or out of range".
      -->
      <Interval>PT1M</Interval>
      <Count>3</Count>
    </RestartOnFailure>
    <Priority>7</Priority>
  </Settings>

  <Actions Context="Author">
    <Exec>
      <!-- EDIT: absolute path to your kb-mcp.exe -->
      <Command>C:\Users\you\.cargo\bin\kb-mcp.exe</Command>
      <Arguments>serve</Arguments>
      <!-- EDIT: directory containing kb-mcp.toml (CWD discovery priority 2) -->
      <WorkingDirectory>C:\Users\you\kb-mcp</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
